# deploy/uran-web.service
[Unit]
Description=Uran dashboard (SvelteKit / adapter-node)
After=network-online.target uran-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/uran-web
EnvironmentFile=/root/uran-web/.env
ExecStart=/usr/bin/node /root/uran-web/build
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
